# files/etc/hotplug.d/usb/23-wwan_modem
#!/bin/sh
#
# Hotplug Script untuk konfigurasi module WWAN
# @lionware <lionware@siltesa.com> - Oktober 2025
# credit mbah wisnu ( remod version by xidz_x )
#

# Cek ACTION
case "$ACTION" in
    add|change|remove) ;;
    *) exit 0 ;;
esac

# Ambil VID
VID=$(echo $PRODUCT | cut -d'/' -f1)

# Cek VID yang didukung
case "$VID" in
    2cb7|1199|8087|413c|1bc7|2c7c|3f0|489|5c6|e8d)
        logger -t wwan_modem "VID:$VID detected - Supported device"
        ;;
    *)
        logger -t wwan_modem "VID:$VID detected - Unsupported device, skipping"
        exit 0
        ;;
esac

# Ambil interface (dengan error handling)
get_iface() {
    uci show network 2>/dev/null | grep "proto='$1'" | head -n1 | cut -d'.' -f2
}

mm_iface=$(get_iface "modemmanager")
atc_iface=$(get_iface "atc")
xmm_iface=$(get_iface "xmm")

# Cek apakah interface ada
interface_exists() {
    uci -q get network."$1" >/dev/null 2>&1
}

# Set port serial
set_port() {
    uci set 3ginfo.@3ginfo[0].device="$1" 2>/dev/null
    uci set atcommands.@atcommands[0].set_port="$1" 2>/dev/null
    uci set modemband.@modemband[0].set_port="$1" 2>/dev/null
    uci set modeminfo.@modeminfo[0].device="$1" 2>/dev/null
    uci set sms_tool_js.@sms_tool_js[0].readport="$1" 2>/dev/null
    uci set sms_tool_js.@sms_tool_js[0].sendport="$1" 2>/dev/null
    uci set sms_tool_js.@sms_tool_js[0].ussdport="$1" 2>/dev/null
    uci set sms_tool_js.@sms_tool_js[0].atport="$1" 2>/dev/null
    uci -q commit 3ginfo
    uci -q commit atcommands
    uci -q commit modemband
    uci -q commit modeminfo
    uci -q commit sms_tool_js
}

# Detect OrangePi Zero3
is_opiz3() {
    if grep -q "OrangePi Zero3" /proc/device-tree/model 2>/dev/null; then
        logger -t wwan_modem "OrangePi Zero3 detected"
        return 0
    else
        logger -t wwan_modem "Not OrangePi Zero3"
        return 1
    fi
}

# Set LED OrangePi Zero3
set_led() {
    if is_opiz3; then
        uci set system.@led[0].dev="$1" 2>/dev/null
        uci -q commit system
        service led restart >/dev/null 2>&1
        logger -t wwan_modem "LED configured for interface: $1"
    fi
}

toggle_mm() {
    local file="/etc/hotplug.d/tty/25-modemmanager-tty"
    [ ! -f "$file" ] && return
    
    if [ "$1" = "enable" ]; then
        sed -i 's/^# \(mm_log "info" .*\)/\1/' "$file"
        sed -i 's/^# \(mm_report_event .*\)/\1/' "$file"
    else
        sed -i 's/^\(mm_log "info" .*\)/# \1/' "$file"
        sed -i 's/^\(mm_report_event .*\)/# \1/' "$file"
    fi
}

set_proto() {
    local proto="$1"
    local iface="$2"
    
    if [ "$proto" = "atc" ] && [ -n "$atc_iface" ]; then
        interface_exists "$atc_iface" && uci set network."$atc_iface".disabled="0"
        if [ -n "$xmm_iface" ] && interface_exists "$xmm_iface"; then
            uci set network."$xmm_iface".disabled="1"
            ifdown "$xmm_iface" 2>/dev/null
        fi
    elif [ "$proto" = "mm" ] && [ -n "$mm_iface" ]; then
        interface_exists "$mm_iface" && uci set network."$mm_iface".disabled="0"
        if [ -n "$xmm_iface" ] && interface_exists "$xmm_iface"; then
            uci set network."$xmm_iface".disabled="1"
            ifdown "$xmm_iface" 2>/dev/null
        fi
    elif [ "$proto" = "xmm" ] && [ -n "$xmm_iface" ]; then
        interface_exists "$xmm_iface" && uci set network."$xmm_iface".disabled="0"
        if [ -n "$mm_iface" ] && interface_exists "$mm_iface"; then
            uci set network."$mm_iface".disabled="1"
            ifdown "$mm_iface" 2>/dev/null
        fi
    fi
    
    uci -q commit network
    
    if [ -n "$iface" ] && interface_exists "$iface"; then
        if [ "$proto" = "atc" ]; then
            ifdown "$iface" 2>/dev/null
            sleep 1
        fi
        ifup "$iface" 2>/dev/null
    fi
}

# Get device name based on VID
get_device_name() {
    case "$1" in
        1199) echo "EM7455" ;;
        2c7c) echo "RM50XX" ;;
        3f0) echo "LT4220" ;;
        2cb7) echo "L850-GL" ;;
        8087) echo "L860-GL" ;;
        413c) echo "DW5821e" ;;
        1bc7) echo "LN940" ;;
        489) echo "Generic" ;;
        5c6) echo "T99W175" ;;
        e8d) echo "FM350-GL" ;;
        *) echo "Unknown" ;;
    esac
}

# Konfigurasi berdasarkan VID
case "$VID" in
    1199) # EM7430/EM7455
        set_port "/dev/ttyUSB1"
        set_led "wwan0"
        toggle_mm enable
        set_proto "mm" "$mm_iface"
        ;;
    2c7c|3f0) # RM5XX, LT4220
        set_port "/dev/ttyUSB2"
        set_led "wwan0"
        toggle_mm enable
        set_proto "mm" "$mm_iface"
        ;;
    2cb7) # L850
        set_port "/dev/ttyACM0"
        set_led "wwan0"
        toggle_mm enable
        set_proto "mm" "$mm_iface"
        ;;
    8087) # L860
        set_port "/dev/ttyACM0"
        set_led "eth1"
        toggle_mm enable
        set_proto "xmm" "$xmm_iface"
        ;;
    413c|1bc7|489) # DW5821e, LN940
        set_port "/dev/ttyUSB0"
        set_led "wwan0"
        toggle_mm enable
        set_proto "mm" "$mm_iface"
        ;;
    5c6) # T99W175
        set_port "/dev/ttyUSB3"
        set_led "wwan0"
        toggle_mm disable
        set_proto "mm" "$mm_iface"
        ;;
    e8d) # FM350-GL
        set_port "/dev/ttyUSB4"
        set_led "eth1"
        toggle_mm disable
        set_proto "atc" "$atc_iface"
        ;;
esac

# Main logic
STATUS_FILE="/tmp/wwan_modem.status"

if [ ! -f "$STATUS_FILE" ] && [ "$ACTION" != "remove" ]; then
    device_path="/sys${DEVPATH}"
    
    if [ -n "$device_path" ]; then
        logger -t wwan_modem "$(get_device_name $VID) detected (VID:$VID)"
        
        # Set device path hanya jika mm_iface ada
        if [ -n "$mm_iface" ] && interface_exists "$mm_iface"; then
            uci set network."$mm_iface".device="$device_path"
            uci -q commit network
        fi
        
        echo "$VID|$(get_device_name $VID)|$(date '+%H:%M:%S')" > "$STATUS_FILE"
        logger -t wwan_modem "$(get_device_name $VID) configured"
    fi
    
elif [ -f "$STATUS_FILE" ] && [ "$ACTION" != "remove" ]; then
    logger -t wwan_modem "Device already configured, skip"
fi

# Cleanup
if [ -f "$STATUS_FILE" ] && [ "$ACTION" = "remove" ]; then
    old_device=$(cut -d'|' -f2 "$STATUS_FILE" 2>/dev/null)
    logger -t wwan_modem "$old_device removed"
    
    if is_opiz3; then
        logger -t wwan_modem "Resetting LED (OPiZ3)"
        uci set system.@led[0].dev="" 2>/dev/null
        uci -q commit system
        service led restart >/dev/null 2>&1
    fi
    
    rm -f "$STATUS_FILE"
fi
